# deb-package/scripts/postinst
#!/bin/sh
set -e

# Define paths
SCRIPT_PATH="/usr/local/bin/cloudflared-tcp.sh"
SERVICE_FILE="/etc/systemd/system/cloudflared-tcp.service"

case "$1" in
    configure)
        # 1. Ensure the helper script is executable if it exists
        if [ -f "$SCRIPT_PATH" ]; then
            chmod +x "$SCRIPT_PATH"
        fi

        # 2. Handle systemd registration and startup
        if [ -d /run/systemd/system ]; then
            # Reload to recognize the new or updated service file
            systemctl daemon-reload >/dev/null 2>&1 || true
            
            # If the service file exists, enable and start it
            if [ -f "$SERVICE_FILE" ]; then
                systemctl enable cloudflared-tcp.service >/dev/null 2>&1 || true
                
                if systemctl is-active --quiet cloudflared-tcp.service; then
                    systemctl restart cloudflared-tcp.service >/dev/null 2>&1 || true
                else
                    systemctl start cloudflared-tcp.service >/dev/null 2>&1 || true
                fi
            fi
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0